# ops/systemd/sheaf.service
[Unit]
Description=Sheaf live server
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/sheaf
EnvironmentFile=/srv/sheaf/.env
Environment=MIX_ENV=dev
ExecStart=/srv/sheaf/bin/serve
Restart=always
RestartSec=3
SyslogIdentifier=sheaf
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=default.target
